# deployment/update-schedule.timer
[Unit]
Description=Update schedule definition

[Time]
OnCalendar=*-*-* 0,12,21:15:00
Persistent=false

[Install]
WantedBy=timers.target